# flamix-autoupdate.path
[Unit]
Description=Watch for Flamix Git Changes

[Path]
PathModified=/opt/flamix/.git/refs/remotes/origin/main
Unit=flamix-autoupdate.service

[Install]
WantedBy=multi-user.target
